# xCAT_and_Masterless_Puppet_Aiming_For_Ideal_Configuration_Management/artifact/puppet/modules/frontend/files/etc/systemd/system/run_puppet.service
[Unit]
Description=Runs Puppet once after booting
After=multi-user.target xcatpostinit1.service

[Service]
Type=forking
ExecStart=/usr/site/rcac/sbin/run_puppet
TimeoutSec=0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
