# xCAT_and_Masterless_Puppet_Aiming_For_Ideal_Configuration_Management/artifact/puppet/modules/frontend/files/etc/cron.d/puppet
30 * * * * root /usr/site/rcac/sbin/run_puppet 2>/dev/null
