# xCAT_and_Masterless_Puppet_Aiming_For_Ideal_Configuration_Management/artifact/puppet/modules/compute/files/etc/cron.d/puppet
30 */4 * * * root /usr/site/rcac/sbin/run_puppet --splay 2>/dev/null
