# Rapid_Deployment_of_Bare-Metal_and_In-Container_HPC_clusters_using_OpenHPC_playbooks/artifact/inception-hpcsyspro18/playbook/roles/master/pbs/files/pbs.service
[Unit]
Description=PBS Professional
After=network-online.target remote-fs.target nss-lookup.target
Wants=network-online.target
DefaultDependencies=true

[Service]
Type=forking
GuessMainPID=no
IgnoreSIGPIPE=no
TimeoutSec=5min
Restart=no
ExecStart=/etc/init.d/pbs start
ExecStop=/etc/init.d/pbs stop

[Install]
WantedBy=multi-user.target
