# cgroups_py_Using_Linux_Control_Groups_and_Systemd_to_Manage_CPU_Time_and_Memory/artifact/cgroups_py.service
[Unit]
Description=cgroups_py CPU throttler for systemd

[Service]
Type=simple
ExecStart=/usr/sbin/cgroups_py
Restart=on-failure

[Install]
WantedBy=multi-user.target
